# dist/init/systemd/alfred.service
[Unit]
Description=Alfred task runner

[Service]
User=alfred
Group=alfred

PIDFile=/var/run/alfred/alfred.pid
ExecStart=/usr/bin/alfred

[Install]
WantedBy=multi-user.target
