# dist/init/systemd/mash.service
[Unit]
Description=Mash task runner

[Service]
User=mash
Group=mash

PIDFile=/var/run/mash/mash.pid
ExecStart=/usr/bin/mash

[Install]
WantedBy=multi-user.target
